# service/heat-api/run
#!/bin/sh

exec runuser -u heat -g heat -- /usr/bin/heat-api 2>&1

