# service/heat-api-cfn/run
#!/bin/sh

exec runuser -u heat -g heat -- /usr/bin/heat-api-cfn 2>&1

