# service/heat-engine/run
#!/bin/sh

exec runuser -u heat -g heat -- /usr/bin/heat-engine 2>&1

